# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── Ollama ────────────────────────────────────────────────────────────────────
[program:ollama]
command=/usr/local/bin/ollama serve
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/ollama.log
stderr_logfile=/var/log/supervisor/ollama_err.log
environment=OLLAMA_MODELS="/root/.ollama/models"

; ── ChromaDB ──────────────────────────────────────────────────────────────────
[program:chromadb]
command=chroma run --host 0.0.0.0 --port 8000 --path /app/chroma_data
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/chromadb.log
stderr_logfile=/var/log/supervisor/chromadb_err.log

; ── Streamlit ─────────────────────────────────────────────────────────────────
[program:streamlit]
command=python3.11 -m streamlit run /app/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
autostart=false
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/streamlit.log
stderr_logfile=/var/log/supervisor/streamlit_err.log
environment=PYTHONPATH="/app"